# etc/cron.d/test
SHELL=/bin/sh
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=""

* * * * * root date >> /tmp/test.log

